# other/logwatcher
#!/sbin/runscript
##!/bin/bash
#
# description: System-V Init-Script for Logwatcher
#
# Get function from functions library
. /etc/init.d/functions.sh

#servicename=${SVCNAME}
servicename=logwatcher

depend() {
	need localmount net
	use dns
	after bootmisc
}

# Start the service logwatcher
start() {
	ebegin "Starting $servicename"
	start-stop-daemon --start --quiet --exec /usr/bin/logwatcher -- -d
	touch /var/lock/$servicename
	eend $?
}

# Restart the service logwatcher
stop() {
	ebegin "Stopping $servicename"
	start-stop-daemon --stop --quiet --exec /usr/bin/logwatcher
	rm -f /var/lock/$servicename
	eend $?
}

# Get services status
status() {
	echo -n "$servicename is ... "
	if [[ -e /var/lock/$servicename ]]
	then
		echo -n "locked and "
	fi
	pidof $servicename 1>/dev/null && echo "running ($(pidof $servicename))" || echo "dead"
	eend $?
}


### main logic ###
#case "$1" in
#  start)
#        start
#        ;;
#  stop)
#        stop
#        ;;
#  status)
#        status
#        ;;
#  restart|reload|condrestart)
#        stop
#        start
#        ;;
#  *)
#        echo $"Usage: $0 {start|stop|restart|reload|status}"
#        exit 1
#esac
#exit 0
